# assets/install/nemo-agent.service
[Unit]
Description=Nemo agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
ExecStart={{INSTALLED_BINARY_PATH}} serve --state-dir {{STATE_DIR}} --host {{HOST}} --port {{PORT}} --dokku-wrapper {{DOKKU_WRAPPER_PATH}}
Restart=on-failure
RestartSec=2
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{STATE_DIR}}
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
